# assets/projects/blackbox_exporter/etc/profile.d/blackbox-exporter.sh
# Global environment variables for blackbox-exporter

export BLACKBOX_EXPORTER_VERSION=0.16.0
export BLACKBOX_EXPORTER_PORT=9115
