# assets/projects/node_exporter/etc/profile.d/node-exporter.sh
# Global environment variables for node-exporter

export NODE_EXPORTER_VERSION=0.18.1
export NODE_EXPORTER_PORT=9100
